# packages/net/hass/files/crontab
*/30 * * * * /usr/lib/hass/sync_state.sh
*/5  * * * * /usr/lib/hass/write_state.sh

